# sync/update-listener.service
[Unit]
Description=PeydX Update Listener
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/peydx/update-listener.py
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
